# os-image/ubuntu/24.04-zeroclaw/zeroclaw.service
[Unit]
Description=ZeroClaw AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zeroclaw gateway
Restart=on-failure
RestartSec=5
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
